# lovelace-kindle-screensaver/rootfs/etc/services.d/node/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Runs the node server
# ==============================================================================
bashio::log.info "Starting Node server..."

bashio::log.debug "$(bashio::api.supervisor GET '/core/api/discovery_info' true)"

export HA_BASE_URL="$(bashio::jq "$(bashio::api.supervisor GET '/core/api/discovery_info' true)" ".base_url")"
export HA_ACCESS_TOKEN="$(bashio::config 'ha_token')"
export HA_SCREENSHOT_URL=$(bashio::config 'screenshot_path')
export CRON_JOB=$(bashio::config 'cron_job')
export RENDERING_TIMEOUT=$(bashio::config 'rendering_timeout')
export RENDERING_DELAY=$(bashio::config 'rendering_delay')
export RENDERING_SCREEN_HEIGHT=$(bashio::config 'rendering_screen_height')
export RENDERING_SCREEN_WIDTH=$(bashio::config 'rendering_screen_width')
export GRAYSCALE_DEPTH=$(bashio::config 'grayscale_depth')
export LANGUAGE=$(bashio::config 'language')
export ROTATION=$(bashio::config 'rotation')
export SCALING=$(bashio::config 'debug')

bashio::log.info "Using base_url: ${HA_BASE_URL}"

cd /app
exec /usr/bin/npm start