# lovelace-dashboard-inkplate-screenshotter/rootfs/etc/services.d/node/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Runs the node server
# ==============================================================================
bashio::log.info "Starting Node server..."

export HA_BASE_URL="$(bashio::config 'ha_base_url')"
export HA_ACCESS_TOKEN="$(bashio::config 'ha_access_token')"
export PAGES=$(bashio::config 'pages')
export LANGUAGE=$(bashio::config 'language')
export CRON_JOB=$(bashio::config 'cron_job')
export RENDERING_TIMEOUT=$(bashio::config 'rendering_timeout')
export RENDERING_DELAY=$(bashio::config 'rendering_delay')

bashio::log.info "Using base_url: ${HA_BASE_URL}"
bashio::log.info "Pages from the list: ${PAGES}"

cd /app
exec /usr/bin/npm start